# files/apt-update.service
[Unit]
Description=Updates apt packages
Wants=apt-update.timer

[Service]
Type=simple
ExecStart=/usr/local/bin/apt-update.sh

[Install]
WantedBy=multi-user.target
